# sv/sxhkd/run
#!/bin/sh
[ -f "$HOME/.config/sxhkd/sxhkdrc-work" ] && exec sxhkd "$HOME/.config/sxhkd/sxhkdrc-work" || exec sxhkd
