# sv/bspwm-event-listener/run
#!/bin/sh

SVDIR=~/service.x sv check bspwm || exit 1
exec bspwm-event-listener
