# sv/de-panel/run
#!/bin/sh

SVDIR=~/service.x sv check bspwm || exit 1
exec de-panel
